# scripts/launchd/dk.mragile.3dpa-intake.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dk.mragile.3dpa-intake</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>__INTAKE_BOOTSTRAP__</string>
    <string>--repo</string>
    <string>__INTAKE_REPO__</string>
    <string>--ios-repo</string>
    <string>__IOS_REPO__</string>
    <string>--expected-origin</string>
    <string>__EXPECTED_ORIGIN__</string>
    <string>--lock</string>
    <string>__INTAKE_LOCK__</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__INSTALL_ROOT__</string>
  <!-- Daily 12:00 local — owner-set 2026-07-10 (supersedes the plan's proposed
       07:45). Keep scripts/intake-runner.config.json "schedule" in sync. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>12</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <!-- ~/Library/Logs always exists; launchd opens these BEFORE the wrapper
       runs, so they must not point into a dir the wrapper itself creates. -->
  <key>StandardOutPath</key>
  <string>/Users/mustafaozturk-macmini/Library/Logs/3dpa-intake.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mustafaozturk-macmini/Library/Logs/3dpa-intake.err.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
