# scripts/launchd/dk.mragile.3dpa-intake-probe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dk.mragile.3dpa-intake-probe</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/mustafaozturk-macmini/dev/Claude/Projects/3dprintassistant/scripts/intake-launchd-probe.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/mustafaozturk-macmini/dev/Claude/Projects/3dprintassistant</string>
  <!-- ~/Library/Logs always exists; launchd opens these BEFORE the script runs,
       so they must not point into the gitignored staging dir the script mkdirs. -->
  <key>StandardOutPath</key>
  <string>/Users/mustafaozturk-macmini/Library/Logs/3dpa-intake-probe.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mustafaozturk-macmini/Library/Logs/3dpa-intake-probe.err.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
